# scripts/com.curio.analyzer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.curio.analyzer</string>
    <key>ProgramArguments</key>
    <array>
        <string>bash</string>
        <string>/Users/netborn/.openclaw/workspace/curio-ui/scripts/analyze-and-report.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/curio-analysis.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/curio-analysis.err</string>
</dict>
</plist>
